# rag-chroma/deploy/smartpi-rag.service
[Unit]
Description=smartpi RAG API
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/smartpi/rag
EnvironmentFile=/home/pi/smartpi/rag/.env
ExecStart=/home/pi/smartpi/rag/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8094
Restart=always
RestartSec=5
StandardOutput=append:/home/pi/smartpi/data/results/rag_api.log
StandardError=append:/home/pi/smartpi/data/results/rag_api.log

[Install]
WantedBy=multi-user.target
