# rag-chroma/deploy/smarttcm-rag-chroma.service
[Unit]
Description=SmartTCM RAG API (Chroma)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/SmartTCM/rag-chroma
EnvironmentFile=/home/pi/SmartTCM/rag-chroma/.env
ExecStart=/home/pi/SmartTCM/rag-chroma/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8095
Restart=always
RestartSec=5
StandardOutput=append:/home/pi/SmartTCM/data/results/rag_chroma_api.log
StandardError=append:/home/pi/SmartTCM/data/results/rag_chroma_api.log

[Install]
WantedBy=multi-user.target
