# edge/smarttcm-voice-agent.service
[Unit]
Description=SmartTCM Raspberry Pi Voice Agent
After=network-online.target smarttcm-edge-bridge.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/SmartTCM
EnvironmentFile=-/home/pi/SmartTCM/config/voice_agent.env
ExecStart=/home/pi/SmartTCM/edge/pi_run_voice_agent.sh
Restart=always
RestartSec=3
StandardOutput=append:/home/pi/SmartTCM/data/results/voice_agent.log
StandardError=append:/home/pi/SmartTCM/data/results/voice_agent.log

[Install]
WantedBy=multi-user.target
