# edge/smarttcm-edge-stream.service
[Unit]
Description=SmartTCM Raspberry Pi MJPEG Detection Stream
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/SmartTCM
EnvironmentFile=-/home/pi/SmartTCM/config/mjpeg_stream.env
ExecStart=/home/pi/SmartTCM/edge/pi_run_mjpeg_stream.sh
Restart=always
RestartSec=3
StandardOutput=append:/home/pi/SmartTCM/data/results/mjpeg_stream.log
StandardError=append:/home/pi/SmartTCM/data/results/mjpeg_stream.log

[Install]
WantedBy=multi-user.target
