# edge/smarttcm-sensor-hub.service
[Unit]
Description=SmartTCM Raspberry Pi Sensor Hub
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/SmartTCM
EnvironmentFile=-/home/pi/SmartTCM/config/sensor_hub.env
ExecStart=/home/pi/SmartTCM/edge/pi_run_sensor_hub.sh
Restart=always
RestartSec=3
StandardOutput=append:/home/pi/SmartTCM/data/results/sensor_hub.log
StandardError=append:/home/pi/SmartTCM/data/results/sensor_hub.log

[Install]
WantedBy=multi-user.target
