# edge/smartpi-sensor-hub.service
[Unit]
Description=smartpi Raspberry Pi Sensor Hub
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/smartpi
EnvironmentFile=-/home/pi/smartpi/config/sensor_hub.env
ExecStart=/home/pi/smartpi/edge/pi_run_sensor_hub.sh
Restart=always
RestartSec=3
StandardOutput=append:/home/pi/smartpi/data/results/sensor_hub.log
StandardError=append:/home/pi/smartpi/data/results/sensor_hub.log

[Install]
WantedBy=multi-user.target
